# webserver/install/init/edudisplej-screenshot-service.service
[Unit]
Description=EduDisplej Screenshot Service
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=edudisplej
Group=edudisplej
WorkingDirectory=/opt/edudisplej
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/edudisplej/.Xauthority"
ExecStart=/bin/bash /opt/edudisplej/init/edudisplej-screenshot-service.sh start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
