# webserver/install/init/edudisplej-api-client.service
[Unit]
Description=EduDisplej API Client
Documentation=https://github.com/nagy-andras-sk/edudisplej
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/edudisplej
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/edudisplej/init/edudisplej-api-client.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
