# webserver/install/init/edudisplej-kiosk.service
[Unit]
Description=EduDisplej Kiosk Mode - Display System
Documentation=https://github.com/nagy-andras-sk/edudisplej
After=network-online.target systemd-user-sessions.service
Wants=network-online.target
Conflicts=getty@tty1.service

[Service]
Type=simple
User=edudisplej
Group=edudisplej
WorkingDirectory=/home/edudisplej
Environment=HOME=/home/edudisplej
Environment=USER=edudisplej
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U

# PAM session for proper logind/seat/VT access
PAMName=login

# TTY configuration for tty1 (main console)
StandardInput=tty
StandardOutput=journal+console
StandardError=journal+console
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# Switch to VT1 before starting
ExecStartPre=/usr/bin/chvt 1

# Run wrapper script to handle initialization and X startup
ExecStart=/opt/edudisplej/init/kiosk-start.sh

# Restart policy
Restart=always
RestartSec=10

# Security and resource limits
# Allow access to graphics and input devices
SupplementaryGroups=video input tty

[Install]
WantedBy=multi-user.target
