# webserver/install/init/edudisplej-kiosk.service
[Unit]
Description=EduDisplej Kiosk Mode
After=network-online.target
Wants=network-online.target
Conflicts=getty@tty1.service

[Service]
Type=simple
User=edudisplej
Group=edudisplej
WorkingDirectory=/home/edudisplej
Environment=HOME=/home/edudisplej
Environment=USER=edudisplej
Environment=DISPLAY=:0
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# Run wrapper script to handle initialization and X startup
ExecStart=/opt/edudisplej/init/kiosk-start.sh

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
