# webserver/install/init/edudisplej-watchdog.service
[Unit]
Description=EduDisplej X Environment Watchdog
After=edudisplej-kiosk.service

[Service]
Type=simple
User=root
ExecStart=/opt/edudisplej/init/edudisplej-watchdog.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
