# client/systemd/edudisplej.service
[Unit]
Description=EduDisplej Kiosk System
Documentation=https://github.com/03Andras/edudisplej
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash /opt/edudisplej/system/edudisplej-init.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
User=root
Environment="DISPLAY=:0"
Environment="HOME=/opt/edudisplej"

[Install]
WantedBy=multi-user.target
